# scripts/openclaw-health-check.timer
[Unit]
Description=OpenClaw Container Health Check Timer
Requires=openclaw-health-check.service

[Timer]
OnBootSec=5s
OnUnitActiveSec=30s
Persistent=true

[Install]
WantedBy=timers.target
